# misc/supervisor.conf
[program:picbed]
command=gunicorn app:app -c picbed.py
directory=/path/to/picbed/src
autostart=true
autorestart=true
redirect_stderr=true
