# misc/supervisor.conf
[program:sapic]
command=gunicorn app:app -c sapicd.py
directory=/path/to/sapic/src
autostart=true
autorestart=true
redirect_stderr=true
